# init/supichi.conf
description "supichi"

start on stopped cloud-final
stop on runlevel [!2345] 
 
respawn
respawn limit 10 5
oom never

env GOOGLE_APPLICATION_CREDENTIALS=/etc/hello/google-api-key.json
export GOOGLE_APPLICATION_CREDENTIALS
 
kill timeout 86400 #If it's given a stop order, this is how long it will take to stop.
limit nofile 1048576 1048576

script
  exec /usr/bin/java -jar /opt/hello/supichi.jar server /etc/hello/supichi.yml
end script
